# aredis/rootfs/etc/services.d/redis/run
#!/usr/bin/with-contenv sh
cd /data
exec redis-server /etc/redis.conf
